# rp64-pwmconfig.service
[Unit]
Description=RockPro64 PWM config service

[Service]
Type=exec
ExecStart=/usr/bin/rp64-pwmconfig

[Install]
WantedBy=basic.target
